# configs/systemd/sudo-pi-backend.service
[Unit]
Description=SUDO-Pi Backend API (FastAPI + Uvicorn)
After=network.target
Wants=network.target

[Service]
Type=exec
User=sudo-pi
Group=sudo-pi
WorkingDirectory=/opt/sudo-pi/backend
EnvironmentFile=/opt/sudo-pi/backend/.env
ExecStart=/opt/sudo-pi/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1 \
    --loop uvloop \
    --no-access-log

Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

# Capabilities for full system management:
#   CAP_NET_ADMIN / CAP_NET_RAW  — network config, raw sockets
#   CAP_SYS_NICE                 — process priority control
#   CAP_DAC_READ_SEARCH          — read any file regardless of owner (file manager)
#   CAP_DAC_OVERRIDE             — write any file regardless of owner (file manager writes, hostapd.conf)
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_NICE CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_NICE CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sudo-pi-backend

[Install]
WantedBy=multi-user.target
